# titan.deployment/airflow/airflow-scheduler.service
[Unit]
Description=airflow-scheduler

[Service]
ExecStart=/bin/bash -c "cd ~/airflow; airflow scheduler"
Restart=always